# .docker/QWorker/crontab
* * * * * /usr/local/bin/php /srv/app/artisan schedule:run >> /dev/stdout
